# systemd/anamnesis-sync.service
[Unit]
Description=anamnesis incremental sync
After=default.target

[Service]
Type=oneshot
WorkingDirectory=%h/projects/anamnesis
Environment=PYTHONPATH=%h/projects/anamnesis
ExecStart=%h/.claude-mem/semantic-env/bin/python -m anamnesis.cli sync

# --- resource limits ---
MemoryMax=2G
MemoryHigh=1500M
CPUQuota=75%
TasksMax=64

# --- scheduling ---
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# --- hardening ---
ProtectSystem=strict
ReadWritePaths=%h/.claude-mem %h/claude-mem-backups
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=default.target
