# systemd/anamnesis-sync.timer
[Unit]
Description=Run anamnesis sync every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
Unit=anamnesis-sync.service

[Install]
WantedBy=timers.target
