# systemd/anamnestic-verify.service
[Unit]
Description=anamnestic daily integrity + drift verification

[Service]
Type=oneshot
WorkingDirectory=%h/projects/anamnestic
Environment=PYTHONPATH=%h/projects/anamnestic
ExecStart=%h/.claude-mem/semantic-env/bin/python -m anamnestic.cli verify

# --- resource limits ---
MemoryMax=512M
CPUQuota=50%
TasksMax=32

Nice=19
IOSchedulingClass=idle

# --- hardening ---
ProtectSystem=strict
ReadWritePaths=%h/.claude-mem
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=default.target
