# systemd/anamnesis-backup.timer
[Unit]
Description=Daily anamnesis backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30min
Unit=anamnesis-backup.service

[Install]
WantedBy=timers.target
