# systemd/anamnesis-backup.service
[Unit]
Description=claude-anamnesis daily backup

[Service]
Type=oneshot
WorkingDirectory=%h/projects/claude-anamnesis
Environment=PYTHONPATH=%h/projects/claude-anamnesis
ExecStart=%h/.claude-mem/semantic-env/bin/python -m anamnesis.cli backup

# --- resource limits ---
MemoryMax=1G
CPUQuota=50%
TasksMax=16

# --- scheduling ---
Nice=19
IOSchedulingClass=idle

# --- hardening ---
ProtectSystem=strict
ReadWritePaths=%h/.claude-mem %h/claude-mem-backups
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=default.target
