# systemd/anamnestic-sync.timer
[Unit]
Description=Run anamnestic sync every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
Unit=anamnestic-sync.service

[Install]
WantedBy=timers.target
