# systemd/claude-mem-worker.service
[Unit]
Description=claude-mem worker (web viewer on :37777 + MCP loopback)
After=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=%h/projects/anamnesis/bin/claude-mem-worker.sh
ExecStop=/bin/bash -c 'nvm_bin=$(ls -d "$HOME"/.nvm/versions/node/v* 2>/dev/null | sort -V | tail -n1)/bin; export PATH="$HOME/.bun/bin:$nvm_bin:$PATH"; npx -y claude-mem stop || true'
Restart=on-failure
RestartSec=10
StandardOutput=append:%h/.claude-mem/logs/worker-stdout.log
StandardError=append:%h/.claude-mem/logs/worker-stderr.log

# --- resource limits ---
MemoryMax=2G
CPUQuota=100%
TasksMax=256

# --- scheduling ---
Nice=5

# --- hardening ---
NoNewPrivileges=yes

[Install]
WantedBy=default.target
