# systemd/mem-ext-backup.timer
[Unit]
Description=Daily claude-mem-ext backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30min
Unit=mem-ext-backup.service

[Install]
WantedBy=timers.target
