# systemd/mem-ext-sync.timer
[Unit]
Description=Run claude-mem-ext sync every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
Unit=mem-ext-sync.service

[Install]
WantedBy=timers.target
